# setup/start/systemd/cloudron-syslog.service
[Unit]
Description=Cloudron Syslog
After=network.target

[Service]
ExecStart=/usr/local/cloudron-syslog/bin/cloudron-syslog --port 2514 --logdir /home/yellowtent/platformdata/logs
WorkingDirectory=/usr/local/cloudron-syslog
Environment="NODE_ENV=production"
Restart=always
User=yellowtent
Group=yellowtent

[Install]
WantedBy=multi-user.target
